# deploy/neo.service
[Unit]
Description=Neo AI Agent
After=network.target

[Service]
Type=simple
User=neo
WorkingDirectory=/opt/neo
ExecStart=/usr/bin/node --import ./deploy/esm-loader.mjs dist/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=NEO_SYSTEMD_UNIT=neo
Environment=NEO_SYSTEMCTL_SCOPE=system
EnvironmentFile=/opt/neo/.env

[Install]
WantedBy=multi-user.target
